# docker/supervisord.conf
; CreditX Ecosystem - Supervisor Configuration
; Manages all microservices within single Hyperlift container

[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ============================================================================
; nginx - Entry point (handles PORT from Hyperlift)
; ============================================================================
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/nginx-stderr.log
stdout_logfile=/app/logs/nginx-stdout.log
priority=100

; ============================================================================
; Frontend - Next.js (standalone mode)
; ============================================================================
[program:frontend]
command=node /app/frontend/server.js
directory=/app/frontend
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/frontend-stderr.log
stdout_logfile=/app/logs/frontend-stdout.log
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0"
priority=200

; ============================================================================
; API Gateway - Express.js
; ============================================================================
[program:api-gateway]
command=node /app/api/dist/index.js
directory=/app/api
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/api-stderr.log
stdout_logfile=/app/logs/api-stdout.log
environment=NODE_ENV="production",PORT="4000"
priority=300

; ============================================================================
; Agent Orchestrator - FastAPI/Python
; ============================================================================
[program:agent]
command=python -m uvicorn app.main:app --host 0.0.0.0 --port 8010 --workers 2
directory=/app/apps/agent
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/agent-stderr.log
stdout_logfile=/app/logs/agent-stdout.log
environment=ENVIRONMENT="production",PORT="8010"
priority=400

; ============================================================================
; CreditX Service - FastAPI/Python (Core compliance service)
; ============================================================================
[program:creditx-service]
command=python -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2
directory=/app/services/creditx-service
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/creditx-stderr.log
stdout_logfile=/app/logs/creditx-stdout.log
environment=ENVIRONMENT="production",PORT="8000"
priority=500

; ============================================================================
; Threat Service - FastAPI/Python
; ============================================================================
[program:threat-service]
command=python -m uvicorn app.main:app --host 0.0.0.0 --port 8001 --workers 1
directory=/app/services/threat-service
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/threat-stderr.log
stdout_logfile=/app/logs/threat-stdout.log
environment=ENVIRONMENT="production",PORT="8001"
priority=600

; ============================================================================
; Guardian Service - FastAPI/Python
; ============================================================================
[program:guardian-service]
command=python -m uvicorn app.main:app --host 0.0.0.0 --port 8002 --workers 1
directory=/app/services/guardian-service
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/guardian-stderr.log
stdout_logfile=/app/logs/guardian-stdout.log
environment=ENVIRONMENT="production",PORT="8002"
priority=700

; ============================================================================
; Group: All Services
; ============================================================================
[group:creditx]
programs=nginx,frontend,api-gateway,agent,creditx-service,threat-service,guardian-service
priority=999
